# third_party/gestures/gestures/lid_touchpad_helper/lid_touchpad_helper.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Lid Monitor (disable/enable TTC on lid close/open)"
author          "chromium-os-dev@chromium.org"

# This daemon maintains a connection to the X server, thus it uses
# non-standard start/stop conditions:
start on login-prompt-visible
stop on stopping ui

respawn
expect fork

env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

exec /usr/sbin/lid_touchpad_helper
